# monitor-test.service
[Unit]
Description=Monitor "test" process and report status

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor-test.sh
